# dist/etc/systemd/system/pycectv.service
[Unit]
Description=pycectv sercice
After=multi-user.target

[Service]
PartOf=graphical-session.target
#Environment="DISPLAY=:0"
#Environment="XAUTHORITY=/home/pi/.Xauthority"
User=pi
ExecStart=/opt/pycectv/bin/startpycectv

[Install]
WantedBy=graphical.target
